# auto-rotate.service
[Unit]
Description=Auto-rotate Yoga Screen

[Service]
Type=forking
ExecStart=/usr/local/bin/auto-rotate

[Install]
WantedBy=multi-user.target
